# src/overlay/rootfs/etc/systemd/system/androidproperties.service
[Unit]
Description=Android properties

[Service]
Type=oneshot
ExecStart=/etc/androidproperties.sh

[Install]
WantedBy=multi-user.target
